# filewatcher/templates/fwr-server.service
[Unit]
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/fwr-server
Restart=always


[Install]
RequiredBy=network.target
